# systemd/open-webui.service
[Unit]
Description=Open WebUI (Ollama chat interface — Docker container)
After=docker.service ollama.service
Requires=docker.service
Wants=ollama.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/docker start open-webui
ExecStop=/usr/bin/docker stop open-webui
Restart=no

[Install]
WantedBy=multi-user.target
