# systemd/open-webui.service
[Unit]
Description=Open WebUI (Ollama chat interface)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=__USER__
ExecStart=__VENV_DIR__/bin/open-webui serve --host 0.0.0.0 --port __PORT__
Restart=always
RestartSec=5
Environment="DATA_DIR=__DATA_DIR__"
Environment="OLLAMA_BASE_URL=http://127.0.0.1:11434"

[Install]
WantedBy=multi-user.target
